# kiosk-reboot.service
[Unit]
Description=Nightly reboot of LJFD kiosk Pi

[Service]
Type=oneshot
ExecStart=/sbin/shutdown -r now
