# usdd-capture-setup.service
[Unit]
Description=LJFD USDD capture pipeline setup (tc358743 -> /dev/video0)
# Order before the switcher, but do NOT depend on multi-user.target: that made a
# cycle (switcher is WantedBy=multi-user.target and After=this), which systemd
# resolves by dropping the switcher's start job. Device readiness is handled by
# Restart=on-failure retries instead.
Before=usdd-switcher.service
StartLimitIntervalSec=0

[Service]
Type=oneshot
RemainAfterExit=yes
# Never let a wedged v4l2/media-ctl call block boot indefinitely.
TimeoutStartSec=90
ExecStart=/usr/local/sbin/usdd-capture-setup.sh
# The USDD may be absent or still latched over-cap at boot; retry until it locks.
# (exit 2 = no lockable signal). Nightly reboot re-runs this cleanly.
Restart=on-failure
RestartSec=15

[Install]
WantedBy=multi-user.target
