# usdd-switcher.service
[Unit]
Description=LJFD USDD alert switcher (dashboard <-> USDD capture overlay)
After=usdd-capture-setup.service kiosk.service
Requires=usdd-capture-setup.service
Wants=kiosk.service

[Service]
Type=simple
# Runs as root: it stops/starts kiosk.service and drives the DRM display via
# kmssink (which hardware-scales the 1080p capture to fill the 4K panel — the
# X software sinks cannot). The kiosk yields the screen for the alert's duration.
RuntimeDirectory=usdd
ExecStart=/usr/local/bin/usdd-switcher.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
